# supervisord.conf
; Runs the two long-lived halves of the gateway in one container:
;   1. the signal-cli JSON-RPC daemon (loopback HTTP on 127.0.0.1:8080), and
;   2. the Python gateway (callback server + inbound SSE listener).
; Both log to the container's stdout/stderr (the /dev/fd/{1,2} trick) so `docker logs`
; / `fly logs` see everything.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

; A control socket + these two sections are what `supervisorctl` connects through; without
; them `supervisorctl ...` fails with "ini file does not include supervisorctl section".
; We use this to stop/start signal-cli for the one-time registration bootstrap.
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:signal-cli]
; SIGNAL_ACCOUNT is injected from the environment (compose .env / Fly secrets). The daemon
; binds HTTP to loopback only — it is never publicly exposed; the gateway reaches it there.
command=signal-cli --config /data/signal-cli -a %(ENV_SIGNAL_ACCOUNT)s daemon --http 127.0.0.1:8080
user=app
environment=HOME="/home/app"
autostart=true
autorestart=true
startsecs=5
; On a fresh machine the account isn't registered yet, so signal-cli exits fast until the
; one-time `register`/`verify` bootstrap is done — give up quietly rather than loop forever;
; the bootstrap explicitly `supervisorctl start`s it afterwards.
startretries=5
stopsignal=TERM
; Let signal-cli flush its double-ratchet state on shutdown before SIGKILL.
stopwaitsecs=900
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:gateway]
command=/app/.venv/bin/python -m bt_signal_gateway
directory=/app
user=app
environment=HOME="/home/app"
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
